# rootfs/etc/services.d/homeseer/run
#!/usr/bin/execlineb -P

trap
{
  term
  {
    echo [run homeseer] TERM received
  }
  hup
  {
    echo [run homeseer] HUP received
  }
}

with-contenv ./hsstart
